# lib/systemd/system/studio-system.service
[Unit]
Description=Studio System CPU Governor Setting

[Service]
Type=oneshot
ExecStart=/lib/systemd/studio

[Install]
WantedBy=multi-user.target
