# usr/lib/systemd/system/ondemand.service.d/studio.conf
[Service]

ExecStart=
ExecStart=/lib/systemd/studio
