# lib/systemd/system/ondemand.service.d/ubuntustudio.conf
[Service]

ExecStart=
ExecStart=/lib/systemd/ubuntustudio
